# dockyard/supervisord.d/dashboard_setup.ini
; one-time launch setup script
;
; WARNING: this script is doing sensitive things like
; modifying tha database ... best is to run as standalone
;
;
[program:dashboard_setup]
directory=%(ENV_DESTDIR)s/dashboard
command=bash setup.sh

; DISABLED
autostart=false

;no automated start of migrations!
;autostart=false
autorestart=false
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=0
; default priority is 999
